# projects/collectiveflow/web/collectiveflow-web.service
[Unit]
Description=CollectiveFlow Web Interface - Horizontal Consensus System
After=network.target

[Service]
Type=simple
User=collectiveflow
Group=collectiveflow
WorkingDirectory=/opt/collectiveflow/web
Environment="COLLECTIVEFLOW_DATA=/opt/collectiveflow/data"
ExecStart=/opt/collectiveflow/web/venv/bin/python /opt/collectiveflow/web/app.py
Restart=on-failure
RestartSec=10

# Security settings - no special privileges
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/collectiveflow/data

[Install]
WantedBy=multi-user.target